# assets/etc/systemd/system/tarnish.service
[Unit]
Description=oxide system service
After=home.mount
StartLimitInterval=600
StartLimitBurst=4
OnFailure=remarkable-fail.service
Conflicts=remux.service draft.service xochitl.service
Before=remarkable-reboot.service remarkable-shutdown.service

[Service]
ExecStart=/opt/bin/tarnish
Restart=on-failure
RestartSec=5
Environment="HOME=/home/root"
Environment="PATH=/opt/bin:/opt/sbin:/opt/usr/bin:/usr/local/bin:/usr/bin:/bin:/usr/local/sbin:/usr/sbin:/sbin"

[Install]
WantedBy=multi-user.target
