# assets/etc/systemd/system/xochitl.service.d/lockscreen_hook.conf
[Service]
Environment=LD_PRELOAD=/usr/lib/lockscreen_hook.so
