# applications/launcherctl/launcherctl-failure.service
[Unit]
Description=Launcherctl failure service
After=home.mount

[Service]
ExecStart=/home/root/.vellum/sbin/launcherctl-failure
